# scripts/server/echoes-server.service
[Unit]
Description=Echoes of the Machine — Consolidated Server (OTA + Startup + Config)
After=network.target
Wants=network.target

[Service]
Type=simple
User=echoes
Group=echoes
WorkingDirectory=/opt/echoes
ExecStart=/opt/echoes/venv/bin/python /opt/echoes/echoes-server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=echoes-server

# The server writes config.json and node_registry.json to /opt/echoes and
# reads /opt/echoes/firmware/ for OTA binaries (populated by build.sh deploy).
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/echoes /var/log/echoes

[Install]
WantedBy=multi-user.target
