# scripts/config_server/echoes-config.service
[Unit]
Description=Echoes of the Machine - Configuration Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/echoes-config
ExecStart=/opt/echoes-config/venv/bin/python /opt/echoes-config/server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=echoes-config

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/echoes-config

[Install]
WantedBy=multi-user.target
